# collector/com.cardtyranno.issuer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cardtyranno.issuer</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/hoody/Documents/Claude/Projects/제휴마케팅 콜렉터/collector/run_issuer.sh</string>
  </array>

  <!-- 매일 오전 7:50 실행 (거주지 캡처 7:30 직후). Hour/Minute 자유 변경 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>50</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/hoody/Documents/Claude/Projects/제휴마케팅 콜렉터/collector/issuer.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hoody/Documents/Claude/Projects/제휴마케팅 콜렉터/collector/issuer.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
